# nvidia/wayland/supervisord.conf
# =============================================================================
# P8 独立觉醒协议 — Hyprland Direct (无 Gamescope)
# Hyprland 作为容器的根 Wayland 合成器，直接驱动 GPU
# =============================================================================

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=5MB
logfile_backups=0
loglevel=info
pidfile=/tmp/supervisord.pid
childlogdir=/tmp
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[include]
files = /etc/supervisor/conf.d/*.conf

[program:dbus-session]
command=bash -c "dbus-daemon --session --nofork --nopidfile --address=unix:path=$XDG_RUNTIME_DIR/dbus-session-bus & sleep 1; ln -sf $XDG_RUNTIME_DIR/dbus-session-bus $XDG_RUNTIME_DIR/bus; wait"
stdout_logfile=/tmp/dbus-session.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=0
redirect_stderr=true
stopasgroup=true
stopsignal=INT
autostart=true
autorestart=true
priority=10

[program:pipewire]
command=bash -c "until [ -S $XDG_RUNTIME_DIR/dbus-session-bus ]; do sleep 0.5; done; /usr/bin/pipewire"
stdout_logfile=/tmp/pipewire.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=0
redirect_stderr=true
stopasgroup=true
stopsignal=INT
autostart=true
autorestart=true
priority=20

[program:pipewire-pulse]
command=bash -c "until [ -S $XDG_RUNTIME_DIR/pipewire-0 ]; do sleep 0.5; done; /usr/bin/pipewire-pulse"
stdout_logfile=/tmp/pipewire-pulse.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=0
redirect_stderr=true
stopasgroup=true
stopsignal=INT
autostart=true
autorestart=true
priority=25

[program:wireplumber]
command=bash -c "until [ -S $XDG_RUNTIME_DIR/pipewire-0 ]; do sleep 0.5; done; /usr/bin/wireplumber"
stdout_logfile=/tmp/wireplumber.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=0
redirect_stderr=true
stopasgroup=true
stopsignal=INT
autostart=true
autorestart=true
priority=30

# =============================================================================
# Hyprland Direct — 独立觉醒
# 不再嵌套在 Gamescope 内，直接作为根 Wayland 合成器
# =============================================================================
[program:hyprland]
command=bash -c "until pw-cli info >/dev/null 2>&1; do sleep 0.5; done; exec /etc/beagle-wind-vnc/bdwind-hyprland.sh"
stdout_logfile=/tmp/hyprland.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=0
redirect_stderr=true
stopasgroup=true
stopsignal=INT
autostart=true
autorestart=true
priority=40

[program:xdg-desktop-portal]
command=bash -c "until ls $XDG_RUNTIME_DIR/wayland-* 1>/dev/null 2>&1; do sleep 0.5; done; WAYLAND_DISPLAY=$(ls $XDG_RUNTIME_DIR/wayland-* 2>/dev/null | grep -v lock | head -1 | xargs basename); export WAYLAND_DISPLAY; export DBUS_SESSION_BUS_ADDRESS=unix:path=$XDG_RUNTIME_DIR/dbus-session-bus; /usr/lib/xdg-desktop-portal"
stdout_logfile=/tmp/xdg-desktop-portal.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=0
redirect_stderr=true
stopasgroup=true
stopsignal=INT
autostart=true
autorestart=true
priority=45

[program:xdg-desktop-portal-hyprland]
command=bash -c "until ls $XDG_RUNTIME_DIR/wayland-* 1>/dev/null 2>&1; do sleep 0.5; done; WAYLAND_DISPLAY=$(ls $XDG_RUNTIME_DIR/wayland-* 2>/dev/null | head -1 | xargs basename); export WAYLAND_DISPLAY; export DBUS_SESSION_BUS_ADDRESS=unix:path=$XDG_RUNTIME_DIR/dbus-session-bus; export HYPRLAND_INSTANCE_SIGNATURE=$(ls -t $XDG_RUNTIME_DIR/hypr/ 2>/dev/null | head -1); export LD_PRELOAD=/opt/gstreamer/patches/gbm_linear_fix.so; /usr/lib/xdg-desktop-portal-hyprland"
stdout_logfile=/tmp/xdg-desktop-portal-hyprland.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=0
redirect_stderr=true
stopasgroup=true
stopsignal=INT
autostart=true
autorestart=true
priority=46

[program:bdwind-gstreamer]
command=bash -c "/etc/beagle-wind-vnc/bdwind-gstreamer.sh"
stdout_logfile=/tmp/bdwind-gstreamer.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=0
redirect_stderr=true
stopasgroup=true
stopsignal=INT
autostart=true
autorestart=true
priority=50

[program:bdwind-gamepad]
command=/etc/beagle-wind-vnc/bdwind-gamepad.sh
stdout_logfile=/tmp/bdwind-gamepad.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=0
redirect_stderr=true
stopasgroup=true
stopsignal=INT
autostart=true
autorestart=true
priority=60

[program:nginx]
command=bash -c "sudo /usr/sbin/nginx -g 'daemon off;'"
stdout_logfile=/tmp/nginx.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=0
redirect_stderr=true
stopasgroup=true
stopsignal=INT
autostart=true
autorestart=true
priority=70
